# Assets/templates/com.clamgui.clamd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<!-- Label for launchd -->
	<key>Label</key>
	<string>com.clamgui.clamd</string>

	<!-- Program to execute -->
	<key>ProgramArguments</key>
	<array>
		<string>/opt/homebrew/sbin/clamd</string>
		<string>--config-file=</string>
		<string>CONFIG_FILE_PATH</string>
	</array>

	<!-- Working directory -->
	<key>WorkingDirectory</key>
	<string>~</string>

	<!-- Run at load (login) -->
	<key>RunAtLoad</key>
	<true/>

	<!-- Keep daemon running -->
	<key>KeepAlive</key>
	<true/>

	<!-- Standard output and error logging -->
	<key>StandardOutPath</key>
	<string>~/Library/Application Support/ClamGUI/clamd-launchd.log</string>
	<key>StandardErrorPath</key>
	<string>~/Library/Application Support/ClamGUI/clamd-launchd-error.log</string>

	<!-- Process type -->
	<key>ProcessType</key>
	<string>Background</string>

	<!-- Nice value (lower priority) -->
	<key>Nice</key>
	<integer>10</integer>

	<!-- Soft resource limits -->
	<key>SoftResourceLimits</key>
	<dict>
		<!-- Maximum number of open files -->
		<key>NumberOfFiles</key>
		<integer>1024</integer>
	</dict>

	<!-- Environment variables -->
	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/opt/homebrew/bin:/opt/homebrew/sbin:/usr/bin:/bin:/usr/sbin:/sbin</string>
	</dict>
</dict>
</plist>
